# roles/bluemind/fixsslconf/templates/etc/systemd/system/bm-nginx.service.d/sslenv.conf
[Service]
# xenial bug: ensure OPENSSL_CONF is defined
Environment=OPENSSL_CONF=/etc/ssl

